# srcs/influxdb/srcs/supervisord.conf
[supervisord]
nodaemon=true
 
[program:influxd]
command= su-exec influxdb /usr/bin/influxd
user=root
autorestart=true 
[program:telegraf]
command= /usr/bin/telegraf 
user=root
autorestart=true