# srcs/nginx/srcs/supervisord.conf
[supervisord]
nodaemon=true
 
[program:sshd]
command=/usr/sbin/sshd -D
user=root
autorestart=true 
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autorestart=true